# scripts/packaging/debian/tine20-tinebase.cron.d
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * *	www-data	/usr/bin/php -f /usr/share/tine20/tine20.php -- --method Tinebase.triggerAsyncEvents
