# ci/built/supervisor/install_tine.ini
[program:install_tine]
command=/install_tine.sh
autostart=true
priority=1023
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0