# ci/dockerimage/supervisor.d/worker.ini
[program:worker]
user=tine20
command=php -d include_path=/etc/tine20/:/usr/share/tine20 /usr/share/tine20/worker.php --config /etc/tine20/actionQueue.ini
autostart=true
autorestart=true
priority=15
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
